# ashi.service
[Unit]
Description=ASHI AI OS Daemon
Documentation=https://github.com/villaex/ashi
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
WorkingDirectory=/home/basitdev/workspace/ashi
ExecStart=/home/basitdev/workspace/ashi/.venv/bin/python functions/ashi_daemon.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/basitdev/workspace/ashi/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ashi

[Install]
WantedBy=default.target
